# files/system/etc/systemd/system/bootc-fetch-apply-updates.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/sleep 30
ExecStart=
ExecStart=/usr/bin/bootc upgrade --quiet
